# hugin.service
[Unit]
Description=Hugin Task Dispatcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/magnus/repos/hugin
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
EnvironmentFile=/home/magnus/repos/hugin/.env
Environment=PATH=/home/magnus/.npm-global/bin:/usr/local/bin:/usr/bin:/bin

# Sandboxing
ProtectSystem=strict
ReadWritePaths=/home/magnus /tmp
NoNewPrivileges=true
PrivateTmp=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=default.target
