# hugin.service
[Unit]
Description=Hugin Task Dispatcher
After=network-online.target munin-memory.service
Wants=network-online.target

[Service]
Type=simple
User=magnus
WorkingDirectory=/home/magnus/hugin
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
EnvironmentFile=/home/magnus/hugin/.env
Environment=PATH=/home/magnus/.npm-global/bin:/usr/local/bin:/usr/bin:/bin

# Sandboxing
ProtectSystem=strict
ReadWritePaths=/home/magnus/hugin /home/magnus/workspace /tmp
NoNewPrivileges=true
PrivateTmp=false

[Install]
WantedBy=multi-user.target
